# src/etc/crontabs/acme
0 0 1 * * acme-renew --all --verbose --retry
